# freshbridge/deploy/freshbridge.service
[Unit]
Description=FreshBridge API Server
After=network.target mysql.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/freshbridge/backend
EnvironmentFile=/opt/freshbridge/backend/.env
ExecStart=/opt/freshbridge/backend/server
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/freshbridge/backend /var/log

[Install]
WantedBy=multi-user.target
